# crontab
# shellcheck disable=all
PATH=/usr/local/bin:/usr/local/sbin:/usr/bin:/usr/sbin:/bin:/sbin
SHELL=/bin/bash

# 袩褉懈 褋褌邪褉褌械 褋懈褋褌械屑褘/泻芯薪褌械泄薪械褉邪
@reboot /bin/bash /app/startup.sh >> /var/log/cron.log 2>&1

# 袩芯写褗械屑 褉械蟹褞屑械 (褉邪蟹 胁 4 褔邪褋邪)
# 小懈屑胁芯谢 锌褉芯褑械薪褌邪 胁 泻褉芯薪褌邪斜械 薪褍卸薪芯 褝泻褉邪薪懈褉芯胁邪褌褜
# sleep 褋 褉邪薪写芯屑薪芯泄 蟹邪写械褉卸泻芯泄 芯褌 1 写芯 5 屑懈薪褍褌 薪褍卸械薪 胁芯-锌械褉胁褘褏 褔褌芯斜褘 褌芯泻械薪 芯斜薪芯胁懈褌褜褋褟 褍褋锌械谢, 邪 胁芯
# 胁褌芯褉褘褏, 褔褌芯斜褘 邪薪褌懈褋锌邪屑 薪械 褋褉邪斜芯褌邪谢
0 */5 * * * /bin/bash -c 'sleep $((60 + RANDOM \% 241))' && /usr/local/bin/python -m hh_applicant_tool update-resumes >> /var/log/cron.log 2>&1

# 袪邪褋褋褘谢邪械屑 芯褌泻谢懈泻懈
0 5-18 * * * /bin/bash -c 'sleep $((60 + RANDOM \% 241))' && /usr/local/bin/python -m hh_applicant_tool apply-vacancies -L /app/letter.txt -f --excluded-filter 'mlm|avon|oriflame|faberlic|amway|herbalife|斜褍泻屑械泻械褉|褋褌邪胁泻懈\\s*薪邪\\s*褋锌芯褉褌|泻邪蟹懈薪芯|谐械屑斜谢懈薪谐|gambling|forex|斜懈薪邪褉薪\\w+芯锌褑懈芯薪|褎懈薪邪薪褋芯胁\\w+锌懈褉邪屑懈写|写褉芯锌褕懈锌锌懈薪谐|懈薪褎芯斜懈蟹薪械褋|蟹邪锌褍褋泻\\s*芯薪谢邪泄薪.?泻褍褉褋|懈薪褎芯锌褉芯写褍泻褌|blockchain|web3|泻褉懈锌褌芯胁邪谢褞褌|谢邪泄褎\\s*泻芯褍褔' >> /var/log/cron.log 2>&1

# 孝械褋褌: 泻邪卸写褘械 5 屑懈薪褍褌 (蟹邪写械褉卸泻邪 写芯 1 屑懈薪褍褌褘, 褔褌芯斜褘 薪械 薪邪泻谢邪写褘胁邪谢懈褋褜)
#*/5 * * * * /bin/bash -c 'sleep $((1 + RANDOM \% 59))' && /usr/local/bin/python -m hh_applicant_tool apply-vacancies >> /var/log/cron.log 2>&1

# 袩褘褌邪械屑褋褟 芯斜薪芯胁懈褌褜 褌芯泻械薪 泻邪卸写褍褞 屑懈薪褍褌褍. 袙 褝褌芯屑 薪械褌 薪械芯斜褏芯写懈屑芯褋褌懈, 屑芯卸薪芯 蟹邪泻芯屑屑械薪褌懈褉芯胁邪褌褜, 褌邪泻 泻邪泻
# 褌芯泻械薪 芯斜薪芯胁谢褟械褌褋褟 邪胁褌芯屑邪褌懈褔械褋泻懈 锌褉懈 谢褞斜芯屑 蟹邪锌褉芯褋械
*/1 * * * * /usr/local/bin/python -m hh_applicant_tool refresh-token

# AI 芯褌胁械褌褘 胁 褔邪褌邪褏 褋 褉邪斜芯褌芯写邪褌械谢褟屑懈 + 蟹邪锌芯谢薪械薪懈械 邪薪泻械褌
0 9-20 * * * /bin/bash -c 'sleep $((60 + RANDOM \% 241))' && /usr/local/bin/python -m hh_applicant_tool reply-employers --use-claude --fill-forms >> /var/log/cron.log 2>&1

# 袪邪褋泻芯屑屑械薪褌懈褉褍泄, 械褋谢懈 褏芯褔械褕褜 褍写邪谢褟褌褜 褔邪褌褘 褋 芯褌泻邪蟹邪屑懈
#0 */1 * * * /bin/bash -c 'sleep $((60 + RANDOM \% 241))' && /usr/local/bin/python -m hh-applicant-tool clear-negotiations -d >> /var/log/cron.log 2>&1

#*/1 * * * * echo "Cron is alive at $(date)" >> /var/log/cron.log 2>&1
